# docker.d/supervisor.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm

[program:php5-fpm-log]
command=tail -f /var/log/php5-fpm.log

[program:nginx]
command=/usr/sbin/nginx
